# files/mithril_signer_preview.in
#!/bin/sh

# PROVIDE: mithril_signer_preview
# REQUIRE: DAEMON cardano_node_preview
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf to enable this service:
#
# mithril_signer_preview_enable:	Set to YES to enable Mithril Signer.
#                               Default: "NO"
#
# mithril_signer_preview_home:		An absolute path to the daemon home directory.
#                               The directory will be created if not exists.
#                               Default: "/var/db/mithril_signer_preview"
#
# mithril_signer_preview_flags:		Any additional command line flags to pass to Mithril.
#                       	Default: "-vvv"
#

. /etc/rc.subr

name=mithril_signer_preview
desc="mithril_signer_preview daemon"
rcvar=mithril_signer_preview_enable
command=/usr/local/bin/mithril-signer

start_cmd="${name}_start"
start_precmd="${name}_prestart"
stop_cmd="${name}_stop"
status_cmd="${name}_status"

extra_commands="status"

load_rc_config $name
: ${mithril_signer_preview_enable:=NO}
: ${mithril_signer_preview_home:="/var/db/mithril_signer_preview"}
: ${mithril_signer_preview_flags:="-vvv"}

pidfile="/var/run/mithril_signer_preview.pid"
logfile="/var/log/mithril_signer_preview.log"
flags="${mithril_signer_preview_flags}"

#export PARTY_ID="pool1xxxxxxxxxxxxxxxxxxxxxxxxxxxxxx" Not needed for Certified operation
export NETWORK="preview"
#export NETWORK_MAGIC="2" not needed any longer.
export AGGREGATOR_ENDPOINT="https://aggregator.api.mithril.network/aggregator"
export RUN_INTERVAL="60000"
export DB_DIRECTORY="/var/db/cardano_node_preview/preview-db"
export CARDANO_NODE_SOCKET_PATH="/var/db/cardano_node_preview/jail/socket/cardano-node.sock"
export CARDANO_CLI_PATH=/usr/local/bin/cardano-cli
export DATA_STORES_DIRECTORY=/var/db/mithril_signer_preview/preview-stores
export KES_SECRET_KEY_PATH=/var/db/mithril_signer_preview/preview-keys/kes.skey
export OPERATIONAL_CERTIFICATE_PATH=/var/db/mithril_signer_preview/preview-keys/pool.cert
export STORE_RETENTION_LIMIT=5

sanity_check()
{
    #if [ ! -f ${mithril_config} ]
    #then
    #    echo "Invalid value for mithril_config: missing file ${mithril_config}"
    #    exit 1
    #fi
    return 0
}

mithril_signer_preview_prestart()
{
    # Create mithril home directory, if does not exist
    if [ ! -d "${mithril_signer_preview_home}/preview-stores" ]; then
        mkdir -p "${mithril_signer_preview_home}/preview-stores"
        mkdir -p "${mithril_signer_preview_home}/preview-keys"
    fi

    sanity_check
}

mithril_signer_preview_start()
{
    check_startmsgs && echo "Starting ${name}."
    cd $mithril_signer_preview_home && /usr/sbin/daemon -p $pidfile -o ${logfile} \
        ${command} ${flags} 2>&1 > /dev/null
}

mithril_signer_preview_stop()
{
    pid=$(check_pidfile "${pidfile}" "${command}")
    if [ -z "${pid}" ]
    then
        echo "${name} is not running"
        return 1
    else
        echo "Stopping ${name}."
        /bin/kill -INT "$pid"
        wait_for_pids "$pid"
    fi
}


mithril_signer_preview_status()
{
    pid=$(check_pidfile "${pidfile}" "${command}")
    if [ -z "${pid}" ]
    then
      echo "${name} is not running"
      return 1
    else
      echo ${name} is running as pid $pid
    fi
}

run_rc_command "$1"
